# dashboard-auto-start.service
[Unit]
Description=RaspberryPi Dashboard Auto-Start Service (WLAN + Hotspot + Docker + Kiosk)
Documentation=https://github.com/HagerJakob/RaspberryPi_Code
After=network.target docker.service
Wants=network.target
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/bin/bash /home/admin/RaspberryPi_Code/start_dashboard.sh
Restart=on-failure
RestartSec=30

# Standard Ausgabe und Fehler zu Journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dashboard-startup

# Environment - Crucial für GUI/Display
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/root"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
Environment="LC_ALL=en_US.UTF-8"

[Install]
WantedBy=multi-user.target
