# backend/deploy/rpi-obd-dashboard.service
[Unit]
Description=Raspberry Pi OBD Dashboard (Wails)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/rpi-obd-dashboard
ExecStart=/opt/rpi-obd-dashboard/obd-dashboard-wails
Restart=on-failure
RestartSec=2
Environment=LOG_FILE_PATH=/var/log/rpi-obd-dashboard/obd_data_log.txt
Environment=LOG_MAX_BYTES=5242880
Environment=SERIAL_PORTS=/dev/serial0,/dev/ttyAMA0,/dev/ttyS0
Environment=TIME_OFFSET_HOURS=1

[Install]
WantedBy=multi-user.target
